# docker/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sisimpur]
command=/start-services.sh
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/sisimpur.log
